# assets/s6-rc.d/postgres-stop/up
#!/usr/bin/execlineb -S0

foreground { s6-echo "Calling /etc/init/postgres-stop" }
/etc/init/postgres-stop
