# images/s6_assets/s6-rc.d/pulpcore-api/run
#!/usr/bin/execlineb -S0
s6-setuidgid pulp
foreground {
  /database/assets/wait_on_database_migrations.sh
}
with-contenv
importas PULP_STATIC_ROOT PULP_STATIC_ROOT
export DJANGO_SETTINGS_MODULE pulpcore.app.settings
importas DJANGO_SETTINGS_MODULE DJANGO_SETTINGS_MODULE
export PULP_SETTINGS /etc/pulp/settings.py
importas PULP_SETTINGS PULP_SETTINGS
export HOME /var/lib/pulp/
importas HOME HOME
importas PULP_GUNICORN_TIMEOUT PULP_GUNICORN_TIMEOUT
importas PULP_API_WORKERS PULP_API_WORKERS
importas PULP_GUNICORN_RELOAD PULP_GUNICORN_RELOAD
importas PULP_API_WORKERS_MAX_REQUESTS PULP_API_WORKERS_MAX_REQUESTS
importas PULP_API_WORKERS_MAX_REQUESTS_JITTER PULP_API_WORKERS_MAX_REQUESTS_JITTER
foreground { s6-echo "Calling /etc/init/pulpcore-api" }
/etc/init/pulpcore-api
