# images/s6_assets/s6-rc.d/pulpcore-api/run
#!/usr/bin/execlineb -S0
s6-setuidgid pulp
foreground {
  /database/assets/wait_on_database_migrations.sh
}
with-contenv
importas PULP_STATIC_ROOT PULP_STATIC_ROOT
export DJANGO_SETTINGS_MODULE pulpcore.app.settings
importas DJANGO_SETTINGS_MODULE DJANGO_SETTINGS_MODULE
export PULP_SETTINGS /etc/pulp/settings.py
importas PULP_SETTINGS PULP_SETTINGS
export HOME /var/lib/pulp/
importas HOME HOME
/usr/local/bin/gunicorn pulpcore.app.wsgi:application --bind "127.0.0.1:24817" --name pulp-api --timeout 90 --access-logfile - --access-logformat "pulp [%({correlation-id}o)s]: %(h)s %(l)s %(u)s %(t)s \"%(r)s\" %(s)s %(b)s \"%(f)s\" \"%(a)s\""
