# images/s6_assets/s6-rc.d/db-fields-key-create/up
#!/usr/bin/execlineb -S0

foreground { s6-echo "Calling /etc/init/db-fields-key-create" }
/etc/init/db-fields-key-create
