# s6_images/assets/s6-rc.d/postgres-prepare/up
#!/usr/bin/execlineb -S0

foreground { s6-echo "Calling /etc/init/postgres-prepare" }
/etc/init/postgres-prepare
