# images/s6_assets/s6-rc.d/nginx/run
#!/usr/bin/execlineb -P
with-contenv
foreground { s6-echo "Calling /etc/init/nginx" }
/etc/init/nginx
