# s6_images/assets/s6-rc.d/pulpcore-worker@2/run
#!/usr/bin/execlineb -P
s6-setuidgid pulp
foreground {
  /database/assets/wait_on_database_migrations.sh
}
export DJANGO_SETTINGS_MODULE pulpcore.app.settings
export PULP_SETTINGS /etc/pulp/settings.py
export PATH /usr/local/bin:/usr/bin/
export HOME /var/lib/pulp/
/usr/local/bin/pulpcore-worker
