# assets/s6-rc.d/postgres-init/up
#!/usr/bin/execlineb -S0

foreground { s6-echo "Calling /etc/init/postgres-init" }
/etc/init/postgres-init
