# images/s6_assets/s6-rc.d/pulpcore-content/run
#!/usr/bin/execlineb -P
s6-setuidgid pulp
foreground {
  /database/assets/wait_on_database_migrations.sh
}
export DJANGO_SETTINGS_MODULE pulpcore.app.settings
export PULP_SETTINGS /etc/pulp/settings.py
export HOME /var/lib/pulp/
/usr/local/bin/gunicorn pulpcore.content:server --bind "127.0.0.1:24816" --name pulp-content --timeout 90 --worker-class "aiohttp.GunicornWebWorker" -w 2 --access-logfile -
