# s6_images/assets/s6-rc.d/pulpcore-api/run
#!/usr/bin/execlineb -S0
s6-setuidgid pulp
foreground {
  /database/assets/wait_on_database_migrations.sh
}
foreground {
  export DJANGO_SETTINGS_MODULE pulpcore.app.settings
  export PULP_CONTENT_ORIGIN localhost
  /usr/local/bin/pulpcore-manager collectstatic --clear --noinput --link
}
export DJANGO_SETTINGS_MODULE pulpcore.app.settings
export PULP_SETTINGS /etc/pulp/settings.py
export HOME /var/lib/pulp/
/usr/local/bin/gunicorn pulpcore.app.wsgi:application --bind "127.0.0.1:24817" --name pulp-api --timeout 90 --access-logfile - --access-logformat "pulp [%({correlation-id}o)s]: %(h)s %(l)s %(u)s %(t)s \"%(r)s\" %(s)s %(b)s \"%(f)s\" \"%(a)s\""
