# assets/s6-rc.d/certs/up
#!/usr/bin/execlineb -S0

foreground { s6-echo "Calling /etc/init/certs" }
/etc/init/certs
